# .config/systemd/user/gui@links-term-code.service.d/override.conf
[Service]
ExecStart=
ExecStart=%h/dev/i3wmonarch/scripts/%I --disable-server --role code-term
