# .config/systemd/user/gui@links-solaar.service.d/override.conf
[Service]
ExecStart=
ExecStart=%h/dev/i3wmonarch/scripts/%I -w hide
