# .config/systemd/user/gui@links-term-htop.service.d/override.conf
[Service]
ExecStart=
ExecStart=%h/dev/i3wmonarch/scripts/%I -e 'htop' -T 'htop' --role 'the_htop' --disable-server
