# .config/systemd/user/gui@links-term-ctop.service.d/override.conf
[Service]
ExecStart=
ExecStart=%h/dev/i3wmonarch/scripts/%I -e 'ctop' -T 'ctop' --role 'the_ctop' --disable-server
